# .config/systemd/user/backup.service
[Unit]
Description=Backup personal files

[Service]
ExecStart=/usr/bin/rclone sync /home/radovan/files gdrive_backup_crypt:
Type=oneshot

